# watchdog/watchdog.service
[Unit]
Description=Aleph-One Self-Healing Watchdog
Documentation=https://github.com/minseongkim-dev/macro-invest-agent-system
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aleph-one/watchdog

# Activate the venv that pip install -r requirements.txt was run in
ExecStart=/opt/aleph-one/watchdog/venv/bin/python watchdog.py

# Restart on any non-zero exit; wait 10 s before retry
Restart=on-failure
RestartSec=10

# Forward stdout/stderr to journald (view with: journalctl -u aleph-watchdog -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aleph-watchdog

# Hard memory cap — watchdog is lightweight
MemoryMax=256M

[Install]
WantedBy=multi-user.target
